# root/etc/services.d/emby/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc /usr/bin/mono \
	/usr/lib/emby/MediaBrowser.Server.Mono.exe \
	-programdata /config -ffmpeg /usr/bin/ffmpeg \
	-ffprobe /usr/bin/ffprobe -restartpath /defaults/restart.sh

